# user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -c 10 -w 50 -f 97
